# bin/obs/com.pai.obs-embed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  PAI Obsidian Embeddings Auto-Update LaunchAgent

  To install:
  1. Copy this file to ~/Library/LaunchAgents/
  2. Replace YOUR_USERNAME with your actual username
  3. Replace /path/to/PAI with the actual path to your PAI installation
  4. Run: launchctl load ~/Library/LaunchAgents/com.pai.obs-embed.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.obs-embed</string>

    <key>Comment</key>
    <string>PAI: Auto-update Obsidian vault embeddings every 30 minutes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/.bun/bin/bun</string>
        <string>run</string>
        <string>/path/to/PAI/bin/obs/obs.ts</string>
        <string>embed</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/YOUR_USERNAME/.bun/bin</string>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/obs-embed.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/obs-embed.err</string>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
